# systemd/noreposts-feed.service
[Unit]
Description=NoReposts ATProto Feed Generator
After=network.target

[Service]
Type=simple
User=vitorpy
Group=http
WorkingDirectory=/var/www/noreposts-feed
ExecStart=/var/www/noreposts-feed/following-no-reposts-feed
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment file
EnvironmentFile=/var/www/noreposts-feed/.env
Environment="RUST_LOG=info"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
RuntimeDirectory=noreposts-feed

[Install]
WantedBy=multi-user.target
